# supervisor.conf
[supervisord]
nodaemon=true

[program:devstack]
user=stack
command=/github/devstack/stack.sh
